# macos/automation/launchd/com.user.opencode-serve.plist
<?xml version="1.0" encoding="UTF-8" ?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<!--
Runs `opencode serve` to keep the opencode server running as a proper macos service.
-->
<dict>
    <key>Label</key>
    <string>com.user.opencode-serve</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/mkuckert/private/dev/env/macos/automation/launchd/opencode-serve.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true />
    <key>KeepAlive</key>
    <true />
    <key>StandardOutPath</key>
    <string>/Users/mkuckert/Library/Logs/opencode/serve.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/mkuckert/Library/Logs/opencode/serve.stderr.log</string>
</dict>
</plist>
